# vitssr.conf
[program:vitssr]
directory=/var/www/vittim.live
command=node server-dist/index.js
numprocs=1
autostart=true
autorestart=true
user=www-data